# supervisord.conf
[supervisord]

[program:initdata]
command=/opt/dse/bin/init-data.sh
